# dockerfiles/ubuntu/supervisord.conf
[supervisord]
nodaemon=true

[program:screen]
command=/bin/bash

