# dockerfiles/rtmp/dist/nginx.supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.rtmp.conf
redirect_stderr=true

