# dockerfiles/apt-cacher-ng/dist/apt-cacher-ng.supervisor.conf
[program:apt-cacher-ng]
command=/usr/sbin/apt-cacher-ng -c /etc/apt-cacher-ng pidfile=/var/run/apt-cacher-ng/pid SocketPath=/var/run/apt-cacher-ng/socket PassThroughPattern=.* Remap-secdeb=security.debian.org ExTreshold=1 ForeGround=1
redirect_stderr=true
user=apt-cacher-ng
group=apt-cacher-ng


