# dockerfiles/apt-cacher-ng/dist/apt-cacher-ng.supervisor.conf
[program:apt-cacher-ng]
command=/usr/sbin/apt-cacher-ng -c /etc/apt-cacher-ng pidfile=/var/run/apt-cacher-ng/pid SocketPath=/var/run/apt-cacher-ng/socket foreground=1
user=apt-cacher-ng


