# dockerfiles/debian/dist/supervisord.conf
[supervisord]
nodaemon=true
redirect_stderr=true

[program:cron]
command=/usr/sbin/cron -f -l
redirect_stderr=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -c5 -n
redirect_stderr=true

[program:sshd]
command=/usr/sbin/sshd -p22 -4 -D
redirect_stderr=true

