# dist/debian.supervisord.conf
[supervisord]
nodaemon=true

[program:screen]
command=/usr/bin/screen -S supervisor

