# dockerfiles/logstash/dist/logstash.supervisor.conf
[program:logstash]
process_name=Shipper
command=/opt/logstash/bin/logstash -f /opt/logstash-agent.conf
user=logstash
startretries=3
redirect_stderr=true
std_out_logfile=NONE
startsecs=3
 
[program:logstash-web]
process_name=Web
command=/opt/logstash/bin/logstash web --backend elasticsearch://127.0.0.1/
user=logstash
startretries=3
redirect_stderr=true
std_out_logfile=NONE
startsecs=3
autostart=false
environment=HOME="/opt/logstash"
 
[group:logstash]
programs=logstash
